# huggingface/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid
childlogdir=/tmp

[program:postgres]
command=/usr/lib/postgresql/16/bin/postgres -D %(ENV_PGDATA)s -h 127.0.0.1 -p 5432 -k /tmp
priority=10
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:redis]
command=/app/huggingface/run-redis.sh
priority=10
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:python-ai]
directory=/app
command=/opt/venv/bin/python -m src.main
priority=20
autostart=true
autorestart=true
startsecs=8
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:java-backend]
directory=/app/backend
command=/opt/java/openjdk/bin/java -Xms128m -Xmx512m -jar /app/backend/app.jar
priority=30
autostart=true
autorestart=true
startsecs=12
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:frontend]
directory=/app/frontend
command=/usr/local/bin/node /app/frontend/server.js
environment=HOSTNAME="127.0.0.1",PORT="3000"
priority=30
autostart=true
autorestart=true
startsecs=5
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true

[program:nginx]
command=/usr/sbin/nginx -c /app/huggingface/nginx.conf -g "daemon off;"
priority=40
autostart=true
autorestart=true
startsecs=3
stopasgroup=true
killasgroup=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
